# platform/deploy/systemd/hotpot-vision@.service
[Unit]
Description=Hotpot Vision Worker (%i)
After=network.target hotpot-hub.service
Wants=hotpot-hub.service

[Service]
Type=simple
User=hotpot
Group=hotpot
WorkingDirectory=/opt/hotpot_smart_ops
Environment=HOTPOT_HUB_URL=http://127.0.0.1:8088
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 edge/stream/vision_worker.py \
  --store-id %i \
  --hub-url ${HOTPOT_HUB_URL} \
  --backend yolo \
  --output-dir demo/data/stores/%i/live \
  --interval 5 --cycles 0
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
